# Library/LaunchAgents/com.user.dotfiles-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.dotfiles-sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/fish</string>
		<string>/Users/yuki.c.watanabe/dotfiles/scripts/auto-sync.fish</string>
	</array>

	<key>WatchPaths</key>
	<array>
		<string>/Users/yuki.c.watanabe/.bashrc</string>
		<string>/Users/yuki.c.watanabe/.bash_profile</string>
		<string>/Users/yuki.c.watanabe/.zshrc</string>
		<string>/Users/yuki.c.watanabe/.zprofile</string>
		<string>/Users/yuki.c.watanabe/.profile</string>
		<string>/Users/yuki.c.watanabe/.cshrc</string>
		<string>/Users/yuki.c.watanabe/.tcshrc</string>
		<string>/Users/yuki.c.watanabe/.tmux.conf</string>
		<string>/Users/yuki.c.watanabe/.gitconfig</string>
		<string>/Users/yuki.c.watanabe/.gitconfig-dena</string>
		<string>/Users/yuki.c.watanabe/.config/fish/config.fish</string>
		<string>/Users/yuki.c.watanabe/.config/fish/fish_plugins</string>
		<string>/Users/yuki.c.watanabe/.config/git/ignore</string>
		<string>/Users/yuki.c.watanabe/.config/starship.toml</string>
		<string>/Users/yuki.c.watanabe/.config/gh/config.yml</string>
		<string>/Users/yuki.c.watanabe/.config/gh/hosts.yml</string>
		<string>/Users/yuki.c.watanabe/.config/karabiner/karabiner.json</string>
		<string>/Users/yuki.c.watanabe/.config/mise/config.toml</string>
		<string>/Users/yuki.c.watanabe/.secrets</string>
		<string>/Users/yuki.c.watanabe/.config/fish/functions</string>
		<string>/Users/yuki.c.watanabe/.ssh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>12</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>ThrottleInterval</key>
	<integer>60</integer>

	<key>RunAtLoad</key>
	<false/>

	<key>StandardOutPath</key>
	<string>/Users/yuki.c.watanabe/Library/Logs/dotfiles-sync.stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/yuki.c.watanabe/Library/Logs/dotfiles-sync.stderr.log</string>
</dict>
</plist>
